# server/deploy/systemd/voice-assistant-server.service
[Unit]
Description=Voice Assistant Server
After=network.target

[Service]
Type=simple
User=voiceassistant
Group=voiceassistant
WorkingDirectory=/opt/javia
Environment="PATH=/opt/javia/venv/bin"
ExecStart=/opt/javia/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --timeout-keep-alive 300
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-assistant-server

[Install]
WantedBy=multi-user.target

